# rootfs/etc/s6/dropbear/run
#!/bin/bash

if [ -f ./setup ]
then
  source setup
fi

INIT="/usr/sbin/dropbear -RFEmwg -p 22"

exec ${INIT}
